# supervisord.conf
[supervisord]
nodaemon=true

[program:apache2]
command=/etc/service/apache2/run

[program:mysqld]
command=/etc/service/mysqld/run

[program:asterisk]
command=/etc/service/asterisk/run
